# backups/syncthing/syncthing-backup.timer
[Unit]
Description=Daily backup for syncthing files

[Timer]
OnCalendar=*-*-* 18:00:00

[Install]
WantedBy=basic.target
