# backups/hd/hd-backup.timer
[Unit]
Description=Daily backup for hard drive files

[Timer]
OnCalendar=*-*-* 19:00:00

[Install]
WantedBy=basic.target
